# components/connectivity/src/agrofelis-unificator/docker/crontab
*/2 * * * * bash -lc /web-pub/cron/test.unificator.running.sh